# system/services/upstart/unoconv.conf
description     "unoconv listener"

start on runlevel [2345]
stop on runlevel [016]
respawn
expect daemon
setuid unoconv
exec unoconv -l &

